# supervisord.conf
[supervisord]
nodaemon=true

[program:counter_api]
command=python3 ./main.py

[program:db]
command=python3 ./listen_for_new_cctv_images_and_store_counts_json.py
