# supervisord.conf
[supervisord]
nodaemon=true

[program:counter_api]
command=python3 ./allinone.py


