# supervisord.conf
[supervisord]
nodaemon=true

[program:counter]
command=python3 /ntk38/scripts/uo-object_counting/app/allinone.py
