# supervisord.conf
[supervisord]
nodaemon=true

[program:counter]
command=python3 ./allinone.py
